# wallpaper.service
[Unit]
Description=Ensuring people lock their screen

[Service]
ExecStart=%h/.wallptool/exec.sh
Type=simple

[Install]
WantedBy=default.target
