# wallpaper.service
[Unit]
Description=meow meow meow meow!!!

[Service]
ExecStart=%h/.kitty_wallp/exec.sh
Type=simple

[Install]
WantedBy=default.target
